# root/opt/postman/etc/supervisord.d/php7-fpm.ini
[program:php7-fpm]
command=/usr/sbin/php-fpm7 --quiet --nodaemonize --fpm-config /etc/php7/php-fpm.conf
directory=/tmp
startsecs=5
startretries=3
exitcodes=0
stopsignal=TERM
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
