# root/opt/postman/etc/supervisord.d/spamd.ini
[program:spamd]
command=/usr/sbin/spamd -m 5 -c --socketpath /run/spamd.sock --socketowner spamassassin --socketgroup mail --socketmode 0600
directory=/tmp
startsecs=5
startretries=3
exitcodes=0
stopsignal=TERM
killasgroup=true
user=mail
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
